# templates/systemd/disk-check.timer.j2
[Unit]
Description=disk-check.service timer

[Timer]
OnBootSec=120s
OnUnitActiveSec=24h
Unit=disk-check.service

[Install]
WantedBy=basic.target
